# postgres_to_es/crontab
*/30 * * * * python3 app/load_data.py >> /logfile
