# postgres_to_es/crontab
*/30 * * * * python3 /usr/src/app/app/load_data.py
